# hier/etc/supervisor/conf.d/unifi.conf
[program:unifi]
command=/usr/bin/java -Xmx1024M -jar /usr/lib/unifi/lib/ace.jar start
autostart=true
autorestart=unexpected
directory=/var/tmp

